# LifeOS/install/LifeOS/TOOLS/com.lifeos.healthsync.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!--
  com.lifeos.healthsync.plist.template - Source template. Install via:
    bun ~/.claude/LIFEOS/TOOLS/InstallHealthSync.ts
  Runs HealthSync pull hourly. The installer substitutes placeholders, writes
  the materialized plist to ~/Library/LaunchAgents/, and bootstraps it with
  launchctl. This template stays clean of principal-specific paths so it ships
  in the public LifeOS release.
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.lifeos.healthsync</string>
  <key>ProgramArguments</key>
  <array>
    <string>{{BUN}}</string>
    <string>{{HOME}}/.claude/LIFEOS/TOOLS/HealthSync.ts</string>
    <string>pull</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>StartInterval</key>
  <integer>3600</integer>
  <key>StandardOutPath</key>
  <string>{{HOME}}/.claude/LIFEOS/MEMORY/OBSERVABILITY/health-sync.log</string>
  <key>StandardErrorPath</key>
  <string>{{HOME}}/.claude/LIFEOS/MEMORY/OBSERVABILITY/health-sync.log</string>
  <key>WorkingDirectory</key>
  <string>{{HOME}}/.claude</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>{{BUN_DIR}}:{{HOME}}/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    <key>HOME</key>
    <string>{{HOME}}</string>
  </dict>
</dict>
</plist>
